# walleon/rootfs/etc/services.d/walleon/run
#!/usr/bin/with-contenv sh
# ==============================================================================
# Starts the server
# ==============================================================================

exec node /opt/walleon/build/bin/www.js